# cron/keepawake.cron
0 * * * * (/var/www/php/cron/testWebServices.sh)
